# sites/ESRF/config_folder/home_host/sync/etc/systemd/system/tango-starter.service.d/mbf.conf
[Unit]
Requires=amc525-firmware.service
After=amc525-firmware.service

[Service]
EnvironmentFile=
EnvironmentFile=/opt/host/.mbf-env


