# sites/ESRF/config_folder/home_host/sync/etc/systemd/system/mbf-backup.service
[Unit]
Description=Daily backup of MBF
After=network.target

[Service]
Type=oneshot

User=dserver
Group=comp

ExecStart=/etc/systemd/system/script-esrf/host/mbf-backup.sh

