# sites/ESRF/INSTALL/common/sync/etc/systemd/system/mbf-backup.timer
[Unit]
Description=Run mbf-backup service every day


[Timer]
OnCalendar=*-*-* 03:14:15
AccuracySecs=1us


[Install]
WantedBy=timers.target
